# site-template/systemd/conserver.service.in
[Unit]
Description=conserver - Multi-user IOC console server
After=network-online.target nss-lookup.target time-sync.target
Wants=network-online.target
ConditionPathExists=@INSTALL_LOCATION@/etc/conserver.cf

[Service]
Type=forking
User=conserver
Group=conserver

RuntimeDirectory=conserver
PIDFile=/run/conserver/conserver.pid

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/run/conserver /var/log/conserver
ReadOnlyPaths=@INSTALL_LOCATION@/etc/conserver.cf

LimitNOFILE=4096
MemoryMax=256M

ExecStart=@INSTALL_LOCATION@/sbin/conserver -d -C @INSTALL_LOCATION@/etc/conserver.cf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
