# site-template/conserver.service.in
[Unit]
Description=conserver - Multi-user IOC console server
After=network.target

[Service]
Type=forking
User=conserver
Group=conserver

# Grant permission to bind privileged ports (< 1024)
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Systemd dynamically creates /run/conserver with proper ownership
RuntimeDirectory=conserver
PIDFile=/run/conserver/conserver.pid

ExecStart=@INSTALL_LOCATION@/sbin/conserver -d -C @INSTALL_LOCATION@/etc/conserver.cf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
